# supervisord.conf
[supervisord]
nodaemon=true

[program:php5-fpm]
command=/usr/bin/php-fpm -F -c /etc/php/fpm-php5.6/ -y /etc/php/fpm-php5.6/php-fpm.conf
autorestart=false
